# network_monitor.service
# /etc/systemd/system/network_monitor.service
[Unit]
Description=Network Monitor
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/local/bin/network_monitor.py
Restart=on-failure
RestartSec=60
StartLimitInterval=300
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
